# supervisor/dockerd.conf
[program:dockerd]
command=bash -c "export DOCKER_HOST=%(ENV_DOCKER_HOST)s && /usr/local/bin/start-dockerd.sh"
autostart=true
autorestart=true
stderr_logfile=/var/log/dockerd.err.log
stdout_logfile=/var/log/dockerd.out.log
